# cli/linux-capture-helper/cliks-capture.service
[Unit]
Description=Cliks privacy-isolated activity capture
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/local/libexec/cliks-capture-helper
# Installer supplies CLIKS_CAPTURE_UID and CLIKS_CAPTURE_CLIENT_EXE in a
# drop-in. The helper accepts only that executable running as that user, owns a
# 0600 socket for the user, and emits only while the user owns an active local
# logind seat. It remains root only to open hotplugged input devices.
User=root
Group=root
RuntimeDirectory=cliks
RuntimeDirectoryMode=0755
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
